# deploy/geo-agent.service
[Unit]
Description=GEO Expert Agent - Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/geo-ai-agent
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/ubuntu/geo-ai-agent/.env
ExecStart=/home/ubuntu/geo-ai-agent/.venv/bin/python -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

