# deploy/geo-agent-frontend.service
[Unit]
Description=GEO Expert Agent - Frontend
After=network.target geo-agent.service
Requires=geo-agent.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/geo-ai-agent/frontend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

